# scripts/edge-node/grani-edge-poll.service.example
[Unit]
Description=GRANI edge poll assignments and apply Xray config
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/grani-edge
ExecStart=/usr/local/bin/grani-edge-poll
