# backend/scripts/edge-node/grani-edge-poll.timer.example
[Unit]
Description=Poll GRANI edge assignments every 90 seconds

[Timer]
Unit=grani-edge-poll.service
OnBootSec=45
OnUnitActiveSec=90
AccuracySec=10
Persistent=true

[Install]
WantedBy=timers.target
